# target/resize_rootfs/resize_rootfs.service
[Unit]
Description=Expand the rootfs partition to full size of the boot device.

[Service]
ExecStart=/usr/bin/resize_rootfs.sh
Type=oneshot
RemainAfterExit=yes
TimeoutStopSec=20
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target

